# sparse/lib/udev/rules.d/100-bluetooth-brcm.rules
# do not edit this file, it will be overwritten on update
# edp17 - remove brcm_patchram_plus and restart rfkill service when rfkill changes status

ACTION=="change", SUBSYSTEM=="rfkill", ATTR{type}=="bluetooth", ATTR{state}=="0", RUN+="/usr/bin/droid/bt-fix.sh"