# sparse/lib/systemd/system/cifs.service
[Unit]
Description=Start cifs module cifs.ko
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/cifsstart.sh