# sparse/lib/systemd/system/dhd.service
[Unit]
Description=Start wifi module dhd.ko
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/dhdstart.sh